# backend/scripts/vote.crontab
0 9 1-7 * 1 /opt/auto_campus_vote_routine.sh >> /var/log/campus_routine.log