# backend/scripts/vote.crontab
* * * * * /opt/auto_campus_vote_routine.sh >> /var/log/campus_routine.log